# install/downtime.service
[Unit]
Description=Downtime Notification Bot

[Service]
ExecStart=/usr/bin/downtime/downtime.run

[Install]
WantedBy=multi-user.target
